# mounts/var-local-media.mount
[Mount]
What=/dev/export/media
Where=/var/local/media
Type=xfs

[Install]
WantedBy=multi-user.target
