# mounts/var-local-pictures.mount
[Mount]
What=/dev/export/pictures
Where=/var/local/pictures
Type=xfs

[Install]
WantedBy=multi-user.target
