# mounts/var-local-spinnaker.mount
[Mount]
What=/dev/export/spinnaker
Where=/var/local/spinnaker
Type=xfs

[Install]
WantedBy=multi-user.target
