# units/acme-sh.timer
[Unit]
Description=Renew acme.sh certificates
ConditionPathExists=/etc/acme.sh/account.conf

[Timer]
OnCalendar=daily
