# units/acme-sh.timer
[Unit]
Description=Renew acme.sh certificates

[Timer]
OnCalendar=daily
RandomizedDelaySec=100min

[Install]
WantedBy=timers.target
