# units/keycloak-backup.timer
[Unit]
Description=Backup Keycloak database

[Timer]
OnCalendar=daily
RandomizedDelaySec=100min

[Install]
WantedBy=timers.target
