# units/var-local-nextcloud.mount
[Mount]
What=/dev/export/nextcloud
Where=/var/local/nextcloud
Type=xfs

[Install]
WantedBy=multi-user.target
