# units/var-local-backups.mount
[Mount]
What=/dev/export/backups
Where=/var/local/backups
Type=xfs

[Install]
WantedBy=multi-user.target
