# units/var-local-marie.mount
[Mount]
What=/dev/export/marie
Where=/var/local/marie
Type=xfs

[Install]
WantedBy=multi-user.target
